# deploy/systemd/fakenet-billing-radius-connector.service
[Unit]
Description=FAKE.NET Billing Radius Connector
PartOf=fakenet-billing-stack.target
Wants=network-online.target redis-server.service postgresql.service freeradius.service
After=network-online.target redis-server.service postgresql.service freeradius.service

[Service]
Type=simple
WorkingDirectory=/opt/fakenet-billing
EnvironmentFile=-/etc/fakenet-billing.env
ExecStart=/usr/bin/node src/radius-connector-service.js
Restart=always
RestartSec=5
KillSignal=SIGTERM
SyslogIdentifier=fakenet-billing-radius-connector

[Install]
WantedBy=multi-user.target
